# net-misc/novacomd/files/novacomd
#!/sbin/runscript

depend() {
        need net
}

start() {
        ebegin "Starting ${SVCNAME}"
        export LD_LIBRARY_PATH="/opt/webos/usr/lib64"

        start-stop-daemon --make-pidfile --pidfile=/run/novacomd.pid --background --start -- /opt/webos/usr/sbin/novacomd  -b

        eend 
}

stop() {
        start-stop-daemon --stop --pidfile=/run/novacomd.pid
}

# vim: ts=4 filetype=gentoo-init-d

