# net-misc/luna-service2/files/ls-hubd-prv
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting ${SVCNAME}"
	export LD_LIBRARY_PATH="/opt/webos/usr/lib64"

	start-stop-daemon --start -- /opt/webos/usr/sbin/ls-hubd --daemon --conf /opt/webos/etc/luna-service2/ls-private.conf --pid-dir=/var/run --syslog

	eend 
}

stop() {
	start-stop-daemon --stop -p /var/run/ls-hubd.private.pid 
}

# vim: ts=4 filetype=gentoo-init-d
